# hypr/.config/udev/99-power-profile.rules
# Auto-switch power profile on AC power change
SUBSYSTEM=="power_supply", ATTR{type}=="Mains", RUN+="/home/netanel/.local/bin/power-profile-auto"
